# icloudpd-rtm.service
[Unit]
Description=Run icloudpd download script nightly for rtm and jcm

[Service]
Type=oneshot
ExecStart=/usr/local/bin/icloudpd-download.sh
